# deploy/bot-g2g-worker.service
[Unit]
Description=BotPasteDon G2G Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botpaste
Group=botpaste
WorkingDirectory=/opt/BotPasteDon
ExecStart=/opt/BotPasteDon/venv/bin/python -m workers.g2g_worker
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
WatchdogSec=120

[Install]
WantedBy=multi-user.target
