# deploy/bot-coordinator.service
[Unit]
Description=BotPasteDon Coordinator (Discord Bot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botpaste
Group=botpaste
WorkingDirectory=/opt/BotPasteDon
ExecStart=/opt/BotPasteDon/venv/bin/python -m coordinator.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
WatchdogSec=120

[Install]
WantedBy=multi-user.target
