# deploy/botpaste.service
[Unit]
Description=BotPasteDon — All G2G Services
After=network-online.target
Wants=network-online.target
# Bring up the supervised watchdog alongside the services. bot-watchdog.service
# is ordered After= this unit and is PartOf= it, so it starts after the services
# and stops when this unit stops.
Wants=bot-watchdog.service

[Service]
Type=forking
ExecStartPre=/bin/bash /opt/BotPasteDon/scripts/stop.sh
ExecStart=/bin/bash /opt/BotPasteDon/scripts/start.sh --clean
ExecStop=/bin/bash /opt/BotPasteDon/scripts/stop.sh
RemainAfterExit=yes
KillMode=none

[Install]
WantedBy=multi-user.target
